# 005. 아톰-모하비-아우룸 협업 파이프라인/deploy/systemd/aurum-deployer.service
[Unit]
Description=Aurum Final Report Deployer
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/caiser77/dgx_workspace/005. 아톰-모하비-아우룸 협업 파이프라인
EnvironmentFile=-/home/caiser77/AI_BASE/pipeline.env
ExecStart=/usr/bin/python3 scripts/aurum_deployer.py --root ${PIPELINE_ROOT} --interval ${AURUM_DEPLOY_INTERVAL}
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
